# services/cloudsurge-job.service
[Unit]
Description=CloudSurge Job to check for cost overruns and running VMs

[Service]
Type=oneshot
ExecStart=%h/.local/bin/cloudsurge-job.sh
